# xray-stats.cron
*/2 * * * * /usr/local/bin/stats-collect
29,59 * * * * /usr/local/bin/stats-shrink
