# xray-stats.cron
PATH=/usr/bin:/bin:/usr/local/bin

*/5 * * * * stats-collect <traffic-data-dir>
01 * * * * stats-shrink <traffic-data-dir>
