# xray-stats.cron
PATH=/usr/bin:/bin:/usr/local/bin

00 * * * * stats-collect <traffic-data-dir>
01 0 * * * stats-shrink <traffic-data-dir>
