# xray-stats.cron
*/5 * * * * /usr/local/bin/stats-collect
01 * * * * /usr/local/bin/stats-shrink
